# SolidRun_HummingBoard_i.MX8/systemd/database.service
[Unit]
Description=Database Service for Embedded Learning Platform
Documentation=https://github.com/yourusername/embedded-learning-platform
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/database-server
Restart=always
RestartSec=10
User=root
Group=root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=database

[Install]
WantedBy=multi-user.target
