# SolidRun_HummingBoard_i.MX8/systemd/embedded-learning-services.target
[Unit]
Description=Embedded Learning Platform Services
Requires=embedded-learning.target
After=embedded-learning.target

[Install]
WantedBy=embedded-learning.target
