# SolidRun_HummingBoard_i.MX8/systemd/embedded-learning.target
[Unit]
Description=Embedded Learning Platform Target
Requires=network.target
After=network.target

[Install]
WantedBy=multi-user.target
