# SolidRun_HummingBoard_i.MX8/systemd/uart.service
[Unit]
Description=UART Communication Service
Documentation=https://github.com/yourusername/embedded-learning-platform
After=systemd-modules-load.service
Wants=systemd-modules-load.service

[Service]
Type=simple
ExecStart=/usr/local/bin/uart-server
Restart=always
RestartSec=5
User=root
Group=root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=uart

[Install]
WantedBy=multi-user.target
